# truffe2/data/supervisord.conf
[unix_http_server]
file=%(here)s/supervisor.sock

[supervisord]
logfile=/var/log/apache2/supervisord.log
pidfile=supervisord.pid

[program:celeryd]
command=celery -A app worker -l info
stdout_logfile=/var/log/apache2/celeryd.log
stderr_logfile=/var/log/apache2/celeryd.log
autostart=true
autorestart=false

[supervisorctl]
serverurl=unix://%(here)s/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface
